# flipkart.service
[Unit]
Description=Uvicorn instance to serve Flipkart RAG Backend
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/flipkart-rag/backend
Environment="PATH=/var/www/flipkart-rag/backend/venv/bin"
EnvironmentFile=/var/www/flipkart-rag/backend/.env
# Using bash to expand the PORT variable from EnvironmentFile. Fallback to 8080.
ExecStart=/bin/bash -c "exec /var/www/flipkart-rag/backend/venv/bin/uvicorn app:app --host 0.0.0.0 --port ${PORT:-8080}"

[Install]
WantedBy=multi-user.target
